# scripts/com.valldun1.ruslan-desktop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.valldun1.ruslan-desktop</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Руслан.app/Contents/MacOS/Руслан</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
        <key>OtherJobEnabled</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>Nice</key>
    <integer>-5</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.valldun1.ruslan-desktop.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.valldun1.ruslan-desktop.stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Applications/Руслан.app/Contents/Resources</string>

    <key>LowPriorityIO</key>
    <false/>

    <key>LaunchEvents</key>
    <dict>
        <key>com.apple.iokit.matching</key>
        <dict/>
    </dict>
</dict>
</plist>
